# formulas/matrix/coturn/coturn-update-external-ip.timer
[Unit]
Description=Run external IP change detection every 5 minutes

[Timer]
OnCalendar=*:00/5

[Install]
WantedBy=timers.target
