# formulas/prometheus/exporter/node_exporter/collectors/systemd-unit-state-collector/systemd-unit-state-collector.service.j2
[Unit]
Description=Collect systemd unit metrics

[Service]
Type=simple
ExecStart={{ exe }} {{ textfile_directory }}/systemd-unit-state.prom
