# formulas/prometheus/exporter/node_exporter/collectors/systemd-unit-state-collector/systemd-unit-state-collector.service.j2
[Unit]
Description=Collect systemd unit metrics
PartOf=dnf-automatic.service
After=dnf-automatic.service

[Service]
Type=oneshot
ExecStart={{ exe }} {{ textfile_directory }}/systemd-unit-state.prom
