# formulas/prometheus/exporter/node_exporter/collectors/systemd-unit-state-collector/systemd-unit-state-collector.timer
[Unit]
Description=Trigger Collect systemd unit metrics

[Timer]
OnBootSec=60 seconds
OnUnitActiveSec=5 minutes

[Install]
WantedBy=timers.target
