# formulas/prometheus/exporter/node_exporter/collectors/needs-restart-collector/needs-restart-collector.timer
[Unit]
Description=Trigger Collect metrics of what services need to be restarted

[Timer]
OnBootSec=0
OnBootSec=60
OnUnitActiveSec=1 hours

[Install]
WantedBy=timers.target
